# backend/etc/supervisor/conf.d/concour-doctora.conf
[program:concour-doctora]
directory=/path/to/concour-doctora/backend
command=/path/to/concour-doctora/backend/venv/bin/gunicorn concour_doctora.wsgi:application --workers 3 --bind 127.0.0.1:8000
user=www-data
autostart=true
autorestart=true
stderr_logfile=/var/log/concour-doctora/err.log
stdout_logfile=/var/log/concour-doctora/out.log
